# modules/_GV/gv_misc/other/server-cron/cron.d/sysstat
*/3 * * * * root sh /home/gvadmin/a4s_httpd_restart_and_alert_if_its_down/a4s_httpd_check.sh
